# services/lowbat/lowbat.service
[Unit]
Description=Low Battery Notification Service

[Service]
ExecStart=%h/.dotfiles/scripts/notify/low-battery-notify.sh
